# ops/docker/services/cms/etc/supervisor.ini
[supervisord]
nodaemon=true
user=root

[unix_http_server]
username=dummy
password=dummy

[supervisorctl]
username=dummy
password=dummy

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:php-fpm]
command=/usr/sbin/php-fpm7 -F -R
autorestart=true
